# serving/rag_app/supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

; ── FastAPI / uvicorn ──────────────────────────────────────
[program:fastapi]
command=sh -c "uvicorn app:app
    --host %(ENV_APP_HOST)s
    --port %(ENV_APP_PORT)s
    --workers %(ENV_APP_WORKERS)s
    --log-level %(ENV_LOG_LEVEL)s
    --proxy-headers
    --forwarded-allow-ips='*'"
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/fastapi.stdout.log
stderr_logfile=/var/log/supervisor/fastapi.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=10

; ── Streamlit ──────────────────────────────────────────────
[program:streamlit]
command=streamlit run streamlit_app.py
    --server.port %(ENV_STREAMLIT_PORT)s
    --server.address 0.0.0.0
    --server.headless true
    --server.enableCORS false
    --server.enableXsrfProtection false
    --browser.gatherUsageStats false
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=8
stdout_logfile=/var/log/supervisor/streamlit.stdout.log
stderr_logfile=/var/log/supervisor/streamlit.stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=20

; ── Process group ──────────────────────────────────────────
[group:nyc-tlc]
programs=fastapi,streamlit